# autostart/thunderbird.desktop
[Desktop Entry]
Type=Application
Name=Thunderbird
Exec=thunderbird
X-GNOME-Autostart-enabled=true

