# systemd/backup.service
[Unit]
Description=Backup home directory

[Service]
Type=oneshot
ExecStart=/home/nemo/Desktop/scripts/shell/bin/backup-linux /dev/sda1
